# etc/supervisor/conf.d/kafka.conf
[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh /etc/kafka/server.properties
autostart=true
autorestart=true